# crontab
0,15,30,45 * * * * /home/bheffron/temperature.sh

