# crontab
0,15,30,45 * * * * /home/*/temperature.sh

